# rootfs/etc/services.d/openvpn/finish
#!/usr/bin/execlineb -S0
redirfd -a 2 /dev/null
s6-svscanctl -t /var/run/s6/services
