# rootfs/etc/services.d/openvpn/run
#!/usr/bin/execlineb -P
redirfd -a 2 /var/log/openvpn/stderr.log
/usr/sbin/openvpn
	--user openvpn --group openvpn
	--log-append /var/log/openvpn/stdout.log
	--config /mnt/openvpn/config.ovpn
	--auth-user-pass /mnt/openvpn/passwd
